# systemd/bambucam.service
[Unit]
Description=BambuCam - Raspberry Pi Camera Streaming
Documentation=https://github.com/fgrfn/bambucam
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=bambucam
Group=bambucam
ExecStart=/opt/bambucam/venv/bin/bambucam
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=15
WorkingDirectory=/opt/bambucam

StandardOutput=journal
StandardError=journal
SyslogIdentifier=bambucam

# Runtime and filesystem hardening. The installer rewrites /opt/bambucam when
# BAMBUCAM_DIR points elsewhere. Supplementary camera groups are loaded from
# the service user's group membership created by the installer.
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
ReadWritePaths=/var/lib/bambucam /etc/bambucam /opt/bambucam

Environment="PYTHONUNBUFFERED=1"
# ProtectHome hides /home from the service. Keep libraries and optional
# per-user state in the writable application data directory instead.
Environment="HOME=/var/lib/bambucam"
EnvironmentFile=-/etc/bambucam/environment

[Install]
WantedBy=multi-user.target
